# etc/systemd/mrlogind@.service
[Unit]
Description=Mrlogind Server

[Service]
ExecStart=-@/usr/sbin/tcpd /usr/sbin/in.mrlogind
StandardInput=socket
