# etc/systemd/mrlogind@.service
[Unit]
Description=Mrlogind Server

[Service]
ExecStart=-/usr/sbin/in.mrlogind
StandardInput=socket
