# etc/systemd/mrlogind.socket
[Unit]
Description=Mrlogind Server Activation Socket

[Socket]
ListenStream=0.0.0.0:541
Accept=true

[Install]
WantedBy=sockets.target
