# etc/systemd/mrshd@.service
[Unit]
Description=Mrshd Server

[Service]
ExecStart=-@/usr/sbin/tcpd /usr/sbin/in.mrshd
StandardInput=socket
KillMode=process
